# .crontab
0 7,12,18   * * * (nice time /home/linuxbog/bin/genererboger ; cd /home/www/cvs.linuxbog.dk/htdocs/; head */sideantal.txt) 
59 * * * * /home/linuxbog/bin/opdaterlinuxbog
