# .crontab
3 0,6,12,18   * * * (nice time /home/linuxbog/cvs/bin/genererboger ; cd /home/www/cvs.linuxbog.dk/htdocs/; head */sideantal.txt) 
59 * * * * /home/linuxbog/cvs/bin/opdaterlinuxbog
